# ruut/etc/profile.d/canberra-gtk-module.sh
export GTK3_MODULES="${GTK3_MODULES:-}${GTK3_MODULES:+:}canberra-gtk-module"
