# ruut/etc/systemd/system/multi-user.target.wants/garuda-pacman-lock.service
[Unit]
Description=Pacman lock cleaner

[Service]
Type=oneshot
ExecStart=rm -f /var/lib/pacman/db.lck

[Install]
WantedBy=multi-user.target
